# quizmo.service
[Unit]
Description=Rust Web Application
After=network.target

[Service]
Type=simple
User=deploy
Environment=PORT=3000
Environment=RUST_LOG=info
WorkingDirectory=/opt/rust-webapp
ExecStart=/opt/rust-webapp/rust-webapp
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target
